# install/subscripts/ccc.cron
15 23 * * * /usr/bin/python $PROJECT_ROOT_DIR/pyInit/addPics.py
0 1 * * * /home/pi/phdisp/scripts/gitpull.sh
0 1 * * * /home/pi/phdisp/scripts/displayRestart.sh
@reboot python /home/pi/phdisp/display/screen_manager.py &
