# install/ccc.cron
15 1 * * * /bin/bash /home/pi/photoDisplayer/clTag.sh
15 4 * * * /bin/bash /home/pi/photoDisplayer/googTag.sh
15 23 * * * /usr/bin/python /home/pi/photoDisplayer/pyInit/addPics.py
